# etc/s6-overlay/s6-rc.d/docker-dns/run
#!/command/execlineb -P
with-contenv
/bin/docker-dns
