# etc/s6-overlay/s6-rc.d/localns/run
#!/command/execlineb -P
with-contenv
cd /etc/zones
/bin/localns
