# etc/s6-overlay/s6-rc.d/coredns/run
#!/command/execlineb -P
/usr/bin/coredns -conf /etc/coredns/Corefile
